# Software/Source/deb_dist/pijuice-1.0/debian/pijuice/lib/systemd/system/pijuice.service
[Unit]
Description=PiJuice status service
After=network.target

[Service]
Type=idle
WorkingDirectory=/var/lib/pijuice/
ExecStart=/usr/bin/pijuice_sys.py
ExecStopPost=/usr/bin/pijuice_sys.py stop
Restart=always

[Install]
WantedBy=multi-user.target
